# deploy/com.alphaomega.secfedclaw.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd job: SECFEDCLAW v0.2 daily run.
     Placeholders __PYTHON__ and __DIR__ are filled in by schedule_install.sh.
     Runs every weekday at 16:35 local time (after the US market close). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphaomega.secfedclaw.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__DIR__/daily.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>16</integer><key>Minute</key><integer>35</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>__DIR__/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__DIR__/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
